# files/auto-shutdown.service
[Unit]
Description="Manages EC2 automatic shutdown based on SSH status"

[Service]
ExecStart=/usr/local/sbin/stop-if-inactive.sh

